# logging/speedfan_logger.cron
# Clean speedfan log
0       0       *       *       *       root    rm /var/log/speedfan.log
